# root-files/opt/flownative/supervisor/etc/conf.d/logrotate.conf
[program:logrotate]
process_name=%(program_name)s
command=%(ENV_LOGROTATE_BASE_PATH)s/sbin/logrotate-cron.sh
autostart=true
autorestart=true
numprocs=1
